# linux.upstart.conf
description "lookout monitoring"
author "dskw"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /srv/lookout/lookout.py
